# rootfs/lib/apk/packages/aria2.list
/etc/config/aria2
/etc/init.d/aria2
/lib/apk/packages/aria2.rusers
/usr/bin/aria2c
