# rootfs/lib/apk/packages/base-config-setting.list
/etc/adv_luci_disabled
/etc/config/base_config
/etc/hotplug.d/iface/15-mgr-ipaddr
/etc/hotplug.d/net/15-mgr-ipaddr
/etc/init.d/auto_usb_mount
/etc/init.d/sta_disable
/etc/init.d/usb_factory_reset
/etc/init.d/wifi-init
/etc/uci-defaults/95-base-config-setting
/lib/preinit/81_mnt_ready
/lib/wifi/wifi-init.sh
/usr/sbin/sta_disable
/usr/sbin/system_reset
/usr/share/base-config-setting/auto_usb_mount.hotplug
/usr/share/base-config-setting/etc/inittab
/usr/share/base-config-setting/etc/shadow
/usr/share/base-config-setting/mount.hotplug
/usr/share/firewall4/output-ttl.include
/usr/share/ucitrack/base-config-setting-network.json
/usr/share/ucitrack/base-config-setting-wireless.json
