# rootfs/lib/apk/packages/base-config-setting-ext4fs.list
/etc/config/fstab
/etc/init.d/disk_ready_done
/etc/init.d/ethtool
/lib/preinit/79_disk_ready
/usr/sbin/diskv
