# rootfs/lib/apk/packages/block-mount.list
/etc/hotplug.d/block/00-media-change
/etc/hotplug.d/block/10-mount
/etc/init.d/fstab
/etc/uci-defaults/10-fstab
/lib/libblkid-tiny.so
/lib/upgrade/keep.d/block-mount
/sbin/block
/usr/sbin/swapoff
/usr/sbin/swapon
