# rootfs/lib/apk/packages/natcapd.list
/etc/config/natcapd
/etc/hotplug.d/iface/20-natcapd
/etc/hotplug.d/iface/21-natcapd.ec2ns
/etc/hotplug.d/net/21-natcapd.ec2ns
/etc/init.d/natcapd
/etc/init.d/natcapd_boot
/etc/uci-defaults/95-natcapd
/usr/sbin/natcapd
/usr/sbin/peer_ctl
/usr/share/natcapd/C_cniplist.set
/usr/share/natcapd/board_mac.lua
/usr/share/natcapd/cacert.pem
/usr/share/natcapd/cniplist.orig.set
/usr/share/natcapd/cniplist.set
/usr/share/natcapd/cnlist_update.sh
/usr/share/natcapd/firewall.include
/usr/share/natcapd/getflix.set
/usr/share/natcapd/gfwlist1.tocn.host.txt
/usr/share/natcapd/gfwlist1.tocn.ip.txt
/usr/share/natcapd/gfwlist_update.sh
/usr/share/natcapd/hkiplist.orig.set
/usr/share/natcapd/local.set
/usr/share/natcapd/natcapd.cone_nat_unused.sh
/usr/share/natcapd/natcapd.config.dubai
/usr/share/natcapd/natcapd.openvpn.sh
/usr/share/natcapd/natcapd.pptpd.sh
/usr/share/natcapd/openvpn/ca.crt
/usr/share/natcapd/openvpn/client.conf
/usr/share/natcapd/openvpn/client.crt
/usr/share/natcapd/openvpn/client.key
/usr/share/natcapd/openvpn/dh2048.pem
/usr/share/natcapd/openvpn/server.crt
/usr/share/natcapd/openvpn/server.key
/usr/share/natcapd/openvpn/ta.key
/usr/share/natcapd/yy.json.sh
/usr/share/natcapd/yy.sn.json.lua
/usr/share/ucitrack/natcap.json
