# rootfs/lib/apk/packages/p910nd.list
/etc/config/p910nd
/etc/hotplug.d/usbmisc/20-p910nd
/etc/init.d/p910nd
/lib/apk/packages/p910nd.rusers
/usr/sbin/p910nd
