# rootfs/lib/apk/packages/urllogger.list
/etc/init.d/urllogger
/usr/sbin/urllogger
/usr/sbin/urllogger-ctl
