# rootfs/lib/apk/packages/usbmuxd.list
/etc/hotplug.d/usb/40-usbmuxd
/etc/init.d/usbmuxd
/lib/upgrade/keep.d/usbmuxd
/usr/sbin/usbmuxd
