# rootfs/lib/apk/packages/usteer.list
/etc/config/usteer
/etc/init.d/usteer
/sbin/usteerd
