# rootfs/lib/apk/packages/uwsgi.list
/etc/init.d/uwsgi
/etc/uwsgi/emperor.ini
/usr/sbin/uwsgi
