# rootfs/lib/apk/packages/xl2tpd.list
/etc/config/xl2tp
/etc/init.d/xl2tpd
/etc/ppp/options.xl2tpd
/etc/xl2tpd/xl2tp-secrets
/etc/xl2tpd/xl2tpd.conf
/lib/netifd/proto/l2tp.sh
/usr/sbin/xl2tpd
/usr/sbin/xl2tpd-control
